# deploy/workhub.service
[Unit]
Description=WorkHub Django Application (Gunicorn)
After=network.target mysql.service

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/workhub
EnvironmentFile=/home/ubuntu/workhub/.env

ExecStart=/home/ubuntu/workhub/venv/bin/gunicorn \
    --config /home/ubuntu/workhub/gunicorn.conf.py \
    workhub.wsgi:application

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
